# launchd/com.veer.agentic-pr.agent-test.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.veer.agentic-pr.agent-test</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_ROOT__/bin/run-poller-launchd.sh</string>
    <string>__CONFIG_PATH__</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_ROOT__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__PROJECT_ROOT__/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>__PROJECT_ROOT__/logs/launchd.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/vsinghthiara/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/vsinghthiara/Library/Python/3.9/bin</string>
    <key>PYTHONPATH</key>
    <string>__PROJECT_ROOT__/src</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
